# systemd/claude-tracker.service
# Template — install.sh génère son propre fichier avec le bon WorkingDirectory.
# Garder comme référence pour une install manuelle. Suppose que le repo est cloné
# dans $HOME/claude-tracker (sinon ajuster les paths %h/claude-tracker/...).

[Unit]
Description=Claude Tracker - dashboard centralisé des instances Claude Code
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
WorkingDirectory=%h/claude-tracker/server
ExecStart=%h/claude-tracker/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8765 --log-level warning
Restart=on-failure
RestartSec=3
StandardOutput=append:%h/claude-tracker/logs/server.log
StandardError=append:%h/claude-tracker/logs/server.log

[Install]
WantedBy=default.target
